# deployment/systemd/healtharchive-annual-search-verify.timer
[Unit]
Description=Run HealthArchive annual search verification (capture once per year when ready)

[Timer]
OnCalendar=*-*-* 03:10:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
