# docs/deployment/systemd/healtharchive-storagebox-sshfs.service
[Unit]
Description=HealthArchive Storage Box mount (sshfs)
After=network-online.target
Wants=network-online.target
ConditionPathExists=/etc/healtharchive/storagebox.env

[Service]
Type=simple

# Defaults (can be overridden in /etc/healtharchive/storagebox.env):
Environment=STORAGEBOX_PORT=23
# Empty means "remote home directory" (recommended default for Storage Boxes).
Environment=STORAGEBOX_REMOTE_PATH=
Environment=STORAGEBOX_MOUNT=/srv/healtharchive/storagebox

# You must set these variables in /etc/healtharchive/storagebox.env:
#   STORAGEBOX_HOST, STORAGEBOX_USER, STORAGEBOX_IDENTITY, STORAGEBOX_UID, STORAGEBOX_GID
#
# Optional (overrides defaults above):
#   STORAGEBOX_PORT, STORAGEBOX_REMOTE_PATH, STORAGEBOX_MOUNT
EnvironmentFile=/etc/healtharchive/storagebox.env

ExecStartPre=/usr/bin/install -d -m 0755 -o root -g root ${STORAGEBOX_MOUNT}
ExecStartPre=/usr/bin/install -d -m 0700 -o root -g root /root/.ssh

# -f keeps sshfs in the foreground so systemd can supervise it.
ExecStart=/usr/bin/sshfs -f -p ${STORAGEBOX_PORT} -o IdentityFile=${STORAGEBOX_IDENTITY} -o allow_other,default_permissions -o uid=${STORAGEBOX_UID},gid=${STORAGEBOX_GID},umask=0027 -o reconnect,ServerAliveInterval=15,ServerAliveCountMax=3,kernel_cache ${STORAGEBOX_USER}@${STORAGEBOX_HOST}:${STORAGEBOX_REMOTE_PATH} ${STORAGEBOX_MOUNT}

ExecStop=/bin/fusermount3 -u ${STORAGEBOX_MOUNT}
Restart=on-failure
RestartSec=5
TimeoutStartSec=30

[Install]
WantedBy=multi-user.target
