# docs/deployment/systemd/healtharchive-change-tracking-dry-run.service
[Unit]
Description=HealthArchive change tracking (dry run)
Wants=network-online.target
After=network-online.target postgresql.service
ConditionPathExists=/etc/healtharchive/backend.env

[Service]
Type=oneshot
User=haadmin
WorkingDirectory=/opt/healtharchive-backend
EnvironmentFile=/etc/healtharchive/backend.env
EnvironmentFile=-/etc/healtharchive/healthchecks.env

# Keep impact low on the single VPS.
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=6

TimeoutStartSec=2h

ExecStart=/opt/healtharchive-backend/scripts/systemd-healthchecks-wrapper.sh --ping-var HEALTHARCHIVE_HC_PING_CHANGE_TRACKING -- /opt/healtharchive-backend/.venv/bin/ha-backend compute-changes --since-days 400 --max-events 200 --dry-run
