# docs/deployment/systemd/healtharchive-crawl-metrics.service
[Unit]
Description=HealthArchive crawl progress metrics (node_exporter textfile)
ConditionPathExists=/etc/healtharchive/backend.env

[Service]
Type=oneshot
User=root
Group=root
EnvironmentFile=/etc/healtharchive/backend.env
ExecStart=/opt/healtharchive-backend/.venv/bin/python3 /opt/healtharchive-backend/scripts/vps-crawl-metrics-textfile.py

