# docs/deployment/systemd/healtharchive-replay-reconcile.timer
[Unit]
Description=Run HealthArchive replay reconciler (capped, safe-by-default)

[Timer]
OnCalendar=*-*-* 02:30:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
