# docs/deployment/systemd/healtharchive-coverage-guardrails.timer
[Unit]
Description=Run HealthArchive coverage regression guardrails (daily)

[Timer]
OnCalendar=*-*-* 04:05:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
