# docs/deployment/systemd/healtharchive-cleanup-automation.timer
[Unit]
Description=Run HealthArchive cleanup automation (weekly)

[Timer]
OnCalendar=Sun *-*-* 04:45:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
